# stage3/files
.
./00-install-packages
./00-install-packages/00-debconf
./00-install-packages/00-run.sh
./00-install-packages/01-run-chroot.sh
./00-install-packages/02-packages
./00-install-packages/02-packages-nr
./00-install-packages/deb-pkgs
./00-install-packages/deb-pkgs/lcars-api_0.0.1+nmu10_armhf.deb
./00-install-packages/deb-pkgs/mqtt2influx_0.0.1_armhf.deb
./00-install-packages/deb-pkgs/ng-unrvl-iotw_0.0.2_armhf.deb
./00-install-packages/influxdb.key
./01-tweaks
./01-tweaks/00-run-chroot.sh
./01-tweaks/00-run.sh
./01-tweaks/rootfs
./01-tweaks/rootfs/boot
./01-tweaks/rootfs/boot/config.txt
./01-tweaks/rootfs/etc
./01-tweaks/rootfs/etc/apache2
./01-tweaks/rootfs/etc/apache2/enable-mods.sh
./01-tweaks/rootfs/etc/apache2/sites-available
./01-tweaks/rootfs/etc/apache2/sites-available/000-default.conf
./01-tweaks/rootfs/etc/dnsmasq.conf
./01-tweaks/rootfs/etc/influxdb
./01-tweaks/rootfs/etc/influxdb/influxdb.conf
./01-tweaks/rootfs/etc/lcars
./01-tweaks/rootfs/etc/localtime
./01-tweaks/rootfs/etc/machine-info
./01-tweaks/rootfs/etc/modules
./01-tweaks/rootfs/etc/mosquitto
./01-tweaks/rootfs/etc/mosquitto/conf.d
./01-tweaks/rootfs/etc/mosquitto/conf.d/unraveltec.conf
./01-tweaks/rootfs/etc/mosquitto/conf.d/websockets.conf
./01-tweaks/rootfs/etc/ntopng
./01-tweaks/rootfs/etc/systemd
./01-tweaks/rootfs/etc/systemd/network
./01-tweaks/rootfs/etc/systemd/network/pan0.netdev
./01-tweaks/rootfs/etc/systemd/network/pan0.network
./01-tweaks/rootfs/etc/systemd/system
./01-tweaks/rootfs/etc/systemd/system/broadwayd.service
./01-tweaks/rootfs/etc/systemd/system/bt-agent.service
./01-tweaks/rootfs/etc/systemd/system/cpustats.service
./01-tweaks/rootfs/etc/systemd/system/dbusage.service
./01-tweaks/rootfs/etc/systemd/system/hwclock-early.service
./01-tweaks/rootfs/etc/systemd/system/iostats2mqtt.service
./01-tweaks/rootfs/etc/systemd/system/nap.service
./01-tweaks/rootfs/etc/systemd/system/rfcomm.service
./01-tweaks/rootfs/etc/systemd/system/setcpugov.service
./01-tweaks/rootfs/etc/systemd/system/vcgencmd-stats.service
./01-tweaks/rootfs/etc/systemd/system/wireshark.service
./01-tweaks/rootfs/etc/telegraf
./01-tweaks/rootfs/etc/telegraf/telegraf.d
./01-tweaks/rootfs/etc/telegraf/telegraf.d/influx-socket-output.conf
./01-tweaks/rootfs/etc/telegraf/telegraf.d/mqtt-consumer.conf
./01-tweaks/rootfs/etc/timezone
./01-tweaks/rootfs/etc/update-motd.d
./01-tweaks/rootfs/etc/update-motd.d/90-ut
./01-tweaks/rootfs/root
./01-tweaks/rootfs/root/.screenrc
./01-tweaks/rootfs/root/.vimrc
./01-tweaks/rootfs/root/.vimrc.local
./01-tweaks/rootfs/root/.zshrc
./01-tweaks/rootfs/root/.zshrc.local
./01-tweaks/rootfs/usr
./01-tweaks/rootfs/usr/lib
./01-tweaks/rootfs/usr/lib/systemd
./01-tweaks/rootfs/usr/lib/systemd/system
./01-tweaks/rootfs/usr/lib/systemd/system/bluetooth.service
./01-tweaks/rootfs/usr/lib/tmpfiles.d
./01-tweaks/rootfs/usr/lib/tmpfiles.d/influxdb.conf
./01-tweaks/rootfs/usr/local
./01-tweaks/rootfs/usr/local/bin
./01-tweaks/rootfs/usr/local/bin/cpustats.py
./01-tweaks/rootfs/usr/local/bin/dbusage.py
./01-tweaks/rootfs/usr/local/bin/hwclock-early.sh
./01-tweaks/rootfs/usr/local/bin/iostats2mqtt.py
./01-tweaks/rootfs/usr/local/bin/mqtt-receive.py
./01-tweaks/rootfs/usr/local/bin/set-cpu-gov-perf.sh
./01-tweaks/rootfs/usr/local/bin/vcgencmd-stats.py
./01-tweaks/rootfs/usr/local/bin/wireless-status-mqtt.py
./01-tweaks/rootfs/var
./01-tweaks/rootfs/var/www
./01-tweaks/rootfs/var/www/api
./01-tweaks/rootfs/var/www/api/system
./01-tweaks/rootfs/var/www/api/system/servicelist.php
./02-raspion
./02-raspion/00-run.sh
./02-raspion/01-run-chroot.sh
./02-raspion/ctraspion
./02-raspion/ntopng
./02-raspion/ntopng/runtimeprefs.json
./EXPORT_IMAGE
./files
./prerun.sh
./Readme.txt
